# config/mycouch.service
[Unit]
Description=MyCouch CouchDB JWT Proxy
After=network.target
Wants=couchdb.service

[Service]
Type=simple
User=jm
WorkingDirectory=/var/www/argw.com/mycouch
Environment="PYTHONPATH=/var/www/argw.com/mycouch/src"
ExecStart=/usr/bin/env uv run uvicorn couchdb_jwt_proxy.main:app --host 0.0.0.0 --port 5985
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
